# provisioning/roles/docker/files/docker_prune.service
[Unit]
Description=Docker Prune

[Service]
User=root
Group=root
Type=oneshot
ExecStart=/usr/local/bin/docker_prune.sh